# deploy/systemd/pokras.target
[Unit]
After=network.target

Wants=pokras.service

[Install]
WantedBy=multi-user.target